# infra/systemd/exchange-template.service
# Template — copy to /etc/systemd/system/exchange-<svc>.service and edit
# the placeholders. Run `infra/systemd/install-units.sh` to render all 8
# at once.
#
# Replace:
#   {{SERVICE}}        — eg "auth-service" (the cmd dir name)
#   {{UNIT}}           — eg "auth"  (the suffix used in unit name; matches
#                        the `for svc in auth wallet ...` loop in Jenkinsfile)
#   {{REPO_DIR}}       — eg "/home/oceanroot/exchange"
#   {{ENV_FILE}}       — eg "/etc/exchange/auth-service.env"
#   {{USER}}           — system user that owns REPO_DIR (eg "oceanroot")

[Unit]
Description=Micro-Exchange · {{SERVICE}}
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{REPO_DIR}}/{{SERVICE}}
EnvironmentFile={{ENV_FILE}}
# Binaries land in a shared bin/ produced by the Jenkinsfile deploy stage.
ExecStart={{REPO_DIR}}/bin/{{SERVICE}}
Restart=always
RestartSec=3
# Crash-loop guard so a broken bundle doesn't fill journald.
StartLimitIntervalSec=60
StartLimitBurst=5
# Hardening — runtime doesn't write outside its work dir.
ProtectSystem=full
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
